# docker/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
priority=10
stdout_logfile=/var/log/supervisor/nginx.log
stderr_logfile=/var/log/supervisor/nginx.log

[program:pnbp-app]
command=/app/start-app.sh
directory=/app
autostart=true
autorestart=true
priority=20
stdout_logfile=/var/log/supervisor/pnbp-app.log
stderr_logfile=/var/log/supervisor/pnbp-app.log
environment=ORACLE_HOME="/opt/oracle",LD_LIBRARY_PATH="/opt/oracle"

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock